# deploy/larpscape-admin.service
[Unit]
Description=Larpscape admin console
After=network.target

[Service]
Type=simple
User=larpscape
Group=larpscape
WorkingDirectory=/opt/larpscape
EnvironmentFile=/etc/larpscape/env
ExecStart=/opt/larpscape/node_modules/.bin/tsx admin/server/index.ts
Restart=always
RestartSec=2
Environment=NODE_ENV=production
# Quote the value: systemd splits unquoted Environment= on spaces, which would
# truncate this to RESTART_CMD=sudo and break the admin panel's deploy-restart.
Environment="RESTART_CMD=sudo /usr/bin/systemctl restart larpscape"

[Install]
WantedBy=multi-user.target
